# internal/systemd/testdata/invalid.service
Description=Invalid Service File

[Unit]
# Missing required fields or invalid syntax
InvalidField=ThisIsNotValid
