# internal/systemd/testdata/simple.service
[Unit]
Description=Testdata Unit

[Service]
Type=simple
ExecStart=tail -f /dev/null
